# task5/super-service.service
[Unit]
Description=Super service

[Service]
ExecStart=/vagrant/super-service.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
